# rules/99-canable.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="can0", RUN+="/bin/systemctl start canable.service"